# mutt/.config/systemd/user/mbsync.service
[Unit]
Description=Mai sync

[Service]
Type=oneshot
ExecStart=/home/chris/.scripts/tools/getmail
